# templates/gadget_web.service
[Unit]
Description=USB Gadget Web Control Interface
After=network.target teslausb-safe-mode.service
Wants=network.target
# Skip in safe mode (3+ reboots in 10 min)
ConditionPathExists=!/run/teslausb-safe-mode
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
# Run as root to bind to port 80 (captive portal) and execute sudo commands
User=root
Group=__TARGET_USER__
WorkingDirectory=__GADGET_DIR__/scripts/web
# Ensure /run/samba exists before namespace setup (smbd is now disabled at boot per #74)
ExecStartPre=/bin/mkdir -p /run/samba
ExecStart=/usr/bin/python3 __GADGET_DIR__/scripts/web/web_control.py
Restart=on-failure
RestartSec=30
TimeoutStopSec=10

# Memory limits for Flask app (prevent runaway memory usage)
MemoryMax=250M
MemoryHigh=200M

# Protect from OOM killer (critical service)
OOMScoreAdjust=-500

# Security settings - relaxed for sudo operations
# Note: Disabled mount namespace isolation so mounts are visible system-wide
MountFlags=shared
# /run/samba is prefixed with `-` so systemd silently ignores it when
# missing on a fresh boot. smbd auto-start is disabled (issue #74); the
# directory only exists after edit_usb.sh starts smbd. Without `-` the
# unit fails with `Failed to set up mount namespacing` / status=226/NAMESPACE.
ReadWritePaths=__GADGET_DIR__ __MNT_DIR__ -/run/samba /tmp

[Install]
WantedBy=multi-user.target
